# uctronics-lcd.service
[Unit]
Description=UCTronics LCD Service

[Service]
ExecStart=/usr/bin/uctronics-lcd-startup

[Install]
WantedBy=multi-user.target
